# rootfs/etc/systemd/system/numLockOnTty.service
[Unit]
Description=numlockOnTty
After=systemd-logind.service

[Service]
ExecStart=/usr/bin/numlockOnTty

[Install]
WantedBy=multi-user.target

