# rootfs/etc/systemd/system/x.service
[Unit]
Description=X

[Service]
ExecStart=/usr/bin/startx
User=tt4
WorkingDirectory=/home/tt4/

[Install]
WantedBy=multi-user.target
